# roles/horde_stats_exporter/templates/horde-exporter.service.j2
[Unit]
Description=AI Horde Prometheus Exporter
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{ horde_stats_exporter_user }}
Group={{ horde_stats_exporter_user }}
WorkingDirectory={{ horde_stats_exporter_install_dir }}
ExecStart={{ horde_stats_exporter_install_dir }}/venv/bin/horde-exporter --config {{ horde_stats_exporter_install_dir }}/exporter_config.yaml
Restart=always
RestartSec=10
StandardOutput=append:{{ horde_stats_exporter_log_file }}
StandardError=append:{{ horde_stats_exporter_log_file }}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ horde_stats_exporter_log_file | dirname }}

[Install]
WantedBy=multi-user.target
